# templates/gitlab-backup.timer
[Unit]
Description=Scheduled GitLab Backup

[Timer]
OnCalendar={{ gitlab_backup_on_calendar }}
Persistent=true

[Install]
WantedBy=multi-user.target
